# deploy/cinderella.service
# Cinderella — the consent-based SimpleX archive bot.
#
# ONE long-lived process (Addendum 1 A2): the SimpleX chat core is embedded
# in-process (native addon) alongside the Fastify admin console. There is NO
# separate simplex-chat daemon and NO exposed network port — the admin console
# binds to 127.0.0.1 and nginx terminates TLS in front of it.
#
# See deploy/RUNBOOK.md for the full install/update/backup procedure.
#
# Runtime identity (Addendum 2 A2.2): runs as the non-root `cinderella` user,
# which owns /var/lib/cinderella (bot chat DB, files, media). Secrets are read by
# systemd as root from the EnvironmentFile (0600 root:root) before privileges are
# dropped, so they never sit in a file the service user can read.

[Unit]
Description=Cinderella — consent-based SimpleX-to-web archive bot
Documentation=https://github.com/saschadaemgen/cinderella
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=cinderella
Group=cinderella
WorkingDirectory=/opt/cinderella
EnvironmentFile=/etc/cinderella/cinderella.env
# XFTP stages downloads in TMPDIR then rename()s them into the files folder; keep
# temp on the SAME filesystem as /var/lib/cinderella or the move fails with EXDEV
# (the default /tmp is a separate tmpfs, isolated further by PrivateTmp below).
Environment=TMPDIR=/var/lib/cinderella/xftp-tmp
ExecStart=/usr/bin/node /opt/cinderella/dist/index.js
Restart=on-failure
RestartSec=5
# The SimpleX core closes its DB on SIGTERM; give it time before SIGKILL.
TimeoutStopSec=30

# --- Runtime data (created + owned by systemd as the service user, 0750) ---
StateDirectory=cinderella
StateDirectoryMode=0750
ReadWritePaths=/var/lib/cinderella

# --- Hardening (Addendum 2 A2.4) ---
NoNewPrivileges=true
# The process needs no capabilities (binds an unprivileged localhost port).
CapabilityBoundingSet=
AmbientCapabilities=
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
RemoveIPC=true
UMask=0077
# TCP to SMP/XFTP relays + PostgreSQL, plus DNS (netlink) and local sockets.
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
# Baseline syscall allowlist for a network service. (V8's JIT needs W|X, so
# MemoryDenyWriteExecute is intentionally NOT set.)
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
